# systemd/bb0-display.service
[Unit]
Description=bb0 OLED status display
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
ExecStart=/usr/local/bin/bb0-display.py
Restart=always
RestartSec=5
User=root
KillMode=mixed
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
