# systemd/beaconbutty-housekeeping.timer
[Unit]
Description=BeaconButty — daily storage housekeeping at 08:00

[Timer]
# Run at 08:00, after the 07:00 beacon report has completed
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
