# systemd/beaconbutty-assets.timer
[Unit]
Description=BeaconButty LAN asset scan — every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
