# systemd/beacon-report.timer
[Unit]
Description=BeaconButty — daily beacon report at 07:00

[Timer]
# Run at 07:00 each morning, after overnight traffic has been analysed
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
