# systemd/beaconbutty-assets.service
[Unit]
Description=BeaconButty LAN asset inventory scan
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/beaconbutty-assets.sh
StandardOutput=append:/var/log/beaconbutty/assets.log
StandardError=append:/var/log/beaconbutty/assets.log
