# systemd/beaconbutty-ja4-history.timer
[Unit]
Description=BeaconButty JA4 fingerprint history — daily 00:30

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
