# systemd/beaconbutty-health.timer
[Unit]
Description=BeaconButty — daily health check at 09:30

[Timer]
OnCalendar=09:30
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
