# systemd/beaconbutty-slow-cadence.timer
[Unit]
Description=BeaconButty — run slow-cadence detector hourly

[Timer]
# RITA imports kick off at *:05 and may take several minutes per outstanding
# day; *:35 leaves a comfortable margin even when catching up at boot.
OnCalendar=*:35
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
