# systemd/beaconbutty-teams-relay-check.service
[Unit]
Description=BeaconButty Teams-relay C2 detector (DragonForce/Backdoor.Turn)
After=zeek.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/beaconbutty-teams-relay-check.py
StandardOutput=append:/var/log/beaconbutty/teams-relay.log
StandardError=append:/var/log/beaconbutty/teams-relay.log
