# systemd/resolv-conf-guard.service
[Unit]
Description=BeaconButty — restore /etc/resolv.conf and mark immutable
DefaultDependencies=no
After=local-fs.target
Before=network-pre.target NetworkManager.service
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/resolv-conf-guard.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
