# systemd/suricata-update.service
[Unit]
Description=BeaconButty — update Suricata rules from Emerging Threats Open
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/suricata-update
# Wait for the Suricata command socket to be ready before reloading rules,
# so this doesn't race on boot when Persistent= triggers immediately.
ExecStartPost=/bin/bash -c 'timeout 60 bash -c "until [ -S /var/run/suricata-command.socket ]; do sleep 2; done"; systemctl reload-or-restart suricata'
