# systemd/tailscale-cert-renew.timer
[Unit]
Description=Monthly Tailscale cert renewal

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
