# systemd/wan-watchdog.timer
[Unit]
Description=BeaconButty — check WAN connectivity every 5 minutes
Requires=wan-watchdog.service

[Timer]
# First check 2 minutes after boot (give DHCP time to settle)
OnBootSec=2min
# Then every 5 minutes
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
